# deploy/systemd/openjobslots-deploy.timer
[Unit]
Description=Poll GitHub and deploy OpenJobSlots when main changes

[Timer]
OnBootSec=90s
OnUnitActiveSec=60s
RandomizedDelaySec=15s
Persistent=true
Unit=openjobslots-deploy.service

[Install]
WantedBy=timers.target
