# deploy/systemd/openjobslots-deploy.timer
[Unit]
Description=Poll GitHub and deploy OpenJobSlots when main changes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
RandomizedDelaySec=60s
Persistent=true
Unit=openjobslots-deploy.service

[Install]
WantedBy=timers.target
